# vendor/sigma/rules/schtask_persistence.yml
title: Scheduled Task Persistence
id: sandbox-010
status: stable
level: medium
description: Detects creation of scheduled tasks for persistence
author: Sandbox Platform
date: 2024/01/15
modified: 2024/01/15
tags:
    - attack.persistence
    - attack.t1053.005
logsource:
    category: process_creation
    product: windows
detection:
    selection:
        process_name|contains:
            - schtasks.exe
        command_line|contains:
            - /create
            - /run
    condition: selection
falsepositives:
    - Legitimate software installation
    - System administration tasks
